# lcars/lcars.service
[Unit]
Description=Screensaver

[Service]
User=pi
Type=simple
ExecStart=/bin/bash /home/pi/.lcars/lcars.sh

[Install]
WantedBy=multi-user.target
