# grassgames/grassgames.service
[Unit]
Description=Screensaver

[Service]
User=pi
Type=simple
ExecStart=/bin/bash /home/pi/.screensavers/grassgames.sh

[Install]
WantedBy=multi-user.target
