# flurry/flurry.service
[Unit]
Description=Screensaver

[Service]
User=pi
Type=simple
ExecStart=/bin/bash /home/pi/.flurry/flurry.sh

[Install]
WantedBy=multi-user.target
