# flurry/screensaver.service
[Unit]
Description=Screensaver

[Service]
User=pi
Type=simple
ExecStart=/bin/bash /home/pi/.screensavers/flurry.sh

[Install]
WantedBy=multi-user.target
